# systemd/brs-system.service
[Unit]
Description=Brainware University BRS Automation API
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=brs
Group=brs
WorkingDirectory=/opt/brs
EnvironmentFile=/opt/brs/.env
ExecStart=/opt/brs/.venv/bin/gunicorn -c /opt/brs/gunicorn.conf.py app:app
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/brs/uploads /opt/brs/output /opt/brs/db /opt/brs/logs
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
CapabilityBoundingSet=

# Resource limits
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
